# deploy/local-vm/ufosighting-yt.timer
[Unit]
Description=Run ufosighting YouTube worker every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
