# deploy/ufosighting-ingest.service
[Unit]
Description=ufosighting.report Sighting-post ingest

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ufosighting
EnvironmentFile=/home/ubuntu/ufosighting/.env
ExecStart=/home/ubuntu/ufosighting/.venv/bin/python ingest.py
