# deploy/ufosighting-sync-full.service
[Unit]
Description=ufosighting.report full moderation sync (30-day window)

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ufosighting
EnvironmentFile=/home/ubuntu/ufosighting/.env
ExecStart=/home/ubuntu/ufosighting/.venv/bin/python sync.py --full
