# deploy/ufosighting-sync.timer
[Unit]
Description=Run moderation sync every 15 minutes

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
